# nginx/crontab
* * * * * (echo -n "$(date +"\%Y-\%m-\%dT\%H:\%M:\%S"): " && sshpass -f ~/sshpass.txt ssh dardic.net '~/vps/update-ip.sh' ) >> ~/update-ip.log 2>&1
